# .config/systemd/user/gsettings-tweaks.service
# ══════════════════════════════════════════════════════════════════════
# Ref: https://wiki.archlinux.org/title/Uniform_look_for_Qt_and_GTK_applications
# Ref: https://wiki.archlinux.org/title/Dark_mode_switching
# Ref: https://wiki.archlinux.org/title/Qt
# Ref: https://wiki.archlinux.org/title/GTK
# Ref: https://wiki.archlinux.org/title/HiDPI
# ══════════════════════════════════════════════════════════════════════
# Ref: https://docs.gtk.org/gtk3/class.Settings.html
# Ref: https://docs.gtk.org/gtk4/class.Settings.html
# ══════════════════════════════════════════════════════════════════════
# To find all similar keys on schema type following command:
# ❯ gsettings list-recursively org.gnome.desktop.interface
# ──────────────────────────────────────────────────────────────────────
# To reset all valuses of keys run following command in terminal:
# ❯ gsettings reset-recursively org.gnome.desktop.interface
#  ──────────────────────────────────────────────────────────────────────
# To get a description and allowed values for a specific key:
# ❯ gsettings describe org.gnome.desktop.interface color-scheme
# ──────────────────────────────────────────────────────────────────────
# To change the value of a specific key (e.g., changing the GTK theme):
# ❯ gsettings set org.gnome.desktop.interface gtk-theme 'adw-gtk3-dark'
# ══════════════════════════════════════════════════════════════════════
[Unit]
Description=Apply GTK/dconf settings (first run only)
After=graphical-session.target

[Service]
Type=oneshot
# ══════════════════════════════════════════════════════════════════════
# GLOBAL COLOR SCHEME (For GTK4 / Libadwaita applications)
# ──────────────────────────────────────────────────────────────────────
# Available values: 'default' (Light), 'prefer-dark' (Dark), 'prefer-light'
#
# NOTE: Modern GNOME apps built with GTK4 and Libadwaita ignore the 'gtk-theme' 
# key entirely. They rely solely on this 'color-scheme' switch to toggle 
# their built-in light or dark styles.
# ══════════════════════════════════════════════════════════════════════
#ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface color-scheme 'prefer-dark'
# ══════════════════════════════════════════════════════════════════════
# GTK THEME (For older GTK3 applications)
# ──────────────────────────────────────────────────────────────────────
# Available values: 'Adwaita', 'Adwaita-dark', 'adw-gtk3', 'adw-gtk3-dark', etc.
#
# NOTE: Legacy GTK3 applications still require a specific theme specified here. 
# Using 'adw-gtk3-dark' is highly recommended as it perfectly mimics the modern 
# Libadwaita (GTK4) look, ensuring visual consistency across the entire system.
# ──────────────────────────────────────────────────────────────────────
# [QT] Additional: QT_QPA_PLATFORMTHEME=gtk3
# ══════════════════════════════════════════════════════════════════════
#ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface gtk-theme 'adw-gtk3-dark'
# [Adwaita | MoreWaita | Numix | Numix-Light | Numix-Circle | Numix-Circle-Light | ...  ]
# [Numix-Adwaita | Numix-Adwaita-Light] - моя кастомная тема
ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface icon-theme 'MoreWaita'
ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface font-name 'Adwaita Sans 11'
#ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface document-font-name 'Adwaita Sans 12'
ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface monospace-font-name 'Adwaita Mono 11'
# ──────────────────────────────────────────────────────────────────────
# Предпочтение, указывающее, должен ли рендеринг шрифтов следовать низкоуровневым настройкам 
# «font-hinting», «font-antialiasing» и «font-rgba-order», или учитывать факторы окружения, такие как разрешение и масштаб экрана.
# Возможные значения:
# [manual] - соблюдение низкоуровневых настроек
# [automatic] - предоставление инструментарию принимать собственные решения
# ──────────────────────────────────────────────────────────────────────
ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface font-rendering 'automatic'
# ──────────────────────────────────────────────────────────────────────
# Тип сглаживания, используемый для отрисовки шрифтов. Возможные значения:
# [none] — без сглаживания; 
# [grayscale] — стандартное чёрно-белое сглаживание;
# [rgba] — субпиксельное сглаживание (только для ЖК-экранов).
# ──────────────────────────────────────────────────────────────────────
ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface font-antialiasing 'grayscale'
# [none | slight | medium | full]
ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface font-hinting 'slight'
# ──────────────────────────────────────────────────────────────────────
# Порядок субпиксельных элементов на ЖК-экране; только если установлен тип сглаживания «rgba»
# [rgba | rgb | bgr | vrgb | vbgr]
# ──────────────────────────────────────────────────────────────────────
ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface font-rgba-order 'rgba'
# [true | false]
ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface enable-animations false
# [true | false]
ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface gtk-enable-primary-paste false
# [true | false]
ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface overlay-scrolling true
# ──────────────────────────────────────────────────────────────────────
# Целое число (Int 32), используемое как коэффициент масштабирования окон для экранов с высоким разрешением. «0» — выбрать автоматически, в зависимости от монитора.
# То же самое что переменная - GDK_SCALE
# [0 | 1 | 2 | 3 ...]
# ──────────────────────────────────────────────────────────────────────
#ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface scaling-factor 0
# ──────────────────────────────────────────────────────────────────────
# Коэффициент масштабирования текста (при изменении коэффициента размер шрифта не изменяется).
# То же самое что переменная - GDK_DPI_SCALE
# [Минимально = 0.5 | По умолчанию = 1.0 | Максимально = 3.0]
# ──────────────────────────────────────────────────────────────────────
#ExecStart=/usr/bin/gsettings set org.gnome.desktop.interface text-scaling-factor 1.0
# [default: appmenu:close]
ExecStart=/usr/bin/gsettings set org.gnome.desktop.wm.preferences button-layout 'appmenu:minimize,maximize,close'
# [default: pc105+inet]
ExecStart=/usr/bin/gsettings set org.gnome.desktop.input-sources xkb-model 'chromebook'
ExecStart=/usr/bin/gsettings set org.gnome.desktop.input-sources sources "[('xkb', 'us'), ('xkb', 'ru')]"
# [default: grp:alt_shift_toggle]
ExecStart=/usr/bin/gsettings set org.gnome.desktop.input-sources xkb-options "['grp:win_space_toggle', 'compose:ralt']"
ExecReload=/usr/bin/gsettings reset-recursively org.gnome.desktop.interface
ExecReload=/usr/bin/gsettings reset-recursively org.gnome.desktop.wm.preferences
ExecReload=/usr/bin/gsettings reset-recursively org.gnome.desktop.input-sources
ExecStartPost=/usr/bin/systemctl --user disable %n
RemainAfterExit=yes

[Install]
WantedBy=graphical-session.target
