# etc/systemd/system/systemd-networkd-wait-online.service.d/10-wait-for-only-one-interface.conf
# Ref: https://wiki.archlinux.org/title/Systemd-networkd#Multiple_interfaces_that_are_not_connected_all_the_time
# ══════════════════════════════════════════════════════════════════════
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --any
